# profile.d/jtier-logs-directory.sh

# jtier apps really need this directory :)
mkdir -p /var/groupon/jtier/logs
