# profile.d/jtier-logs-directory.sh

# jtier apps really need this directory :)
mkdir ~/logs
mkdir -p ~/groupon/jtier/logs
